# demos/Shopware/rootfs/etc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stderr
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0

[program:shopware-empty-cache]
command=/usr/local/bin/shopware-empty-cache/shopware-empty-cache.sh
autorestart=unexpected
exitcodes=0
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
